# ansible/roles/undercloud/files/bootif-fix.service
[Unit]
Description=Automated fix for incorrect iPXE BOOTIF 

[Service]
Type=simple
ExecStart=/usr/bin/bootif-fix

[Install]
WantedBy=multi-user.target
